# deploy/cron/iiot-backup-verify.cron.example
# Replace /data/iiot-platform/cloud/deploy with the real Cloud deploy directory before installing.
# Runs a weekly PostgreSQL restore verification at 03:30 every Sunday in the server's local time zone.
30 3 * * 0 cd /data/iiot-platform/cloud/deploy && ./scripts/postgres-verify-backup.sh >> ./backups/postgres/cron-backup-verify.log 2>&1
