# deploy/cron/iiot-backup.cron.example
# Replace /srv/iiot-cloud/deploy with the real deploy directory before installing.
# Runs a daily PostgreSQL backup at 02:30 in the server's local time zone.
30 2 * * * cd /srv/iiot-cloud/deploy && ./scripts/postgres-backup.sh >> ./backups/postgres/cron-backup.log 2>&1
